# devops/crontab
0 */1 * * * /bin/bash /app/devops/processor.sh 
